# cron/crontab
# All Money Back Me Home — Scheduled Tasks
# Install: crontab cron/crontab

# Hourly scanner — every hour at :05
5 * * * * cd /opt/ambmh && docker compose -f docker-compose.prod.yml exec -T api python -m src.jobs.hourly_scanner >> /var/log/ambmh/scanner.log 2>&1

# Daily market report — every day at 08:00 UTC (16:00 Taiwan)
0 8 * * * cd /opt/ambmh && docker compose -f docker-compose.prod.yml exec -T api python -m src.jobs.daily_report >> /var/log/ambmh/daily.log 2>&1

# Institutional report collection + AI digest — every 6 hours (00/06/12/18 UTC)
0 */6 * * * cd /opt/ambmh && docker compose -f docker-compose.prod.yml exec -T telegram-bot python -m market_monitor.report_collector >> /var/log/ambmh/digest.log 2>&1

# Polymarket Phase 1 — every 5 minutes
# Logging handled internally; status at data/reports/polymarket_pipeline_status.json
*/5 * * * * /opt/ambmh/scripts/polymarket_pipeline.sh
