# cron/crontab
# All Money Back Me Home — Scheduled Tasks
# Install: crontab cron/crontab

# Hourly scanner — every hour at :05
5 * * * * cd /opt/ambmh && docker compose -f docker-compose.prod.yml exec -T api python -m src.jobs.hourly_scanner >> /var/log/ambmh/scanner.log 2>&1

# Daily market report — every day at 08:00 UTC (16:00 Taiwan)
# R133: 之前是 'docker compose exec -T api python -m src.jobs.daily_report'
# 但 apps/api/src/jobs/daily_report.py 不存在 (只有 hourly_scanner.py +
# run_backtest.py). silent fail (ModuleNotFoundError 進 daily.log 但無 alert),
# 不知 silent 多久 — 可能從這個 cron line 寫入起就一直壞掉. 改用 host
# 直接 invoke scripts/daily_report.py (跟 polymarket 同 pattern).
0 8 * * * cd /opt/ambmh && /usr/bin/python3 scripts/daily_report.py >> /var/log/ambmh/daily.log 2>&1

# Institutional report collection + AI digest — every 6 hours (00/06/12/18 UTC)
0 */6 * * * cd /opt/ambmh && docker compose -f docker-compose.prod.yml exec -T telegram-bot python -m market_monitor.report_collector >> /var/log/ambmh/digest.log 2>&1

# Polymarket Phase 1 — every 5 minutes
# Logging handled internally; status at data/reports/polymarket_pipeline_status.json
*/5 * * * * /opt/ambmh/scripts/polymarket_pipeline.sh

# R138: Daily log rotation — keep 14 天 cron + script logs, 防 silent disk-full.
# polymarket 每 5 min 跑一次 (288/天), 加 hourly_scanner + daily/digest, 一週累積
# 可達 GB 級。VPS root partition 滿了 = 所有 service 寫不進 log → cascading
# silent failure。每天 04:30 UTC (低流量) 砍 14 天前的 .log 檔。
30 4 * * * find /var/log/ambmh -name "*.log" -mtime +14 -delete 2>/dev/null
